# overlay/etc/init.d/coginfer
#!/sbin/openrc-run
description="CognitiveOS Wide Model Inference Server"

depend() {
    need localmount
    keyword -stop
}

command=/usr/local/bin/coginfer
command_args="--backend cgo --models /cognitiveos/models"
command_background=true
pidfile=/run/coginfer.pid
output_log=/cognitiveos/logs/coginfer.log
error_log=${output_log}
