# overlay/etc/init.d/cograw
#!/sbin/openrc-run
description="CognitiveOS Raw Model Guardrail Server"

depend() {
    need localmount cpm-boot-deps
    keyword -stop
}

command=/usr/local/bin/cograw
command_args="--model /cognitiveos/models/raw/raw-model.gguf"
command_background=true
pidfile=/run/cograw.pid
output_log=/cognitiveos/logs/cograw.log
error_log=${output_log}
