# overlay/etc/init.d/cpm-boot-deps
#!/sbin/openrc-run

description="CognitiveOS Boot-Stage Dependency Collector"

depend() {
    need localmount
    # This is a one-shot service
    keyword -stop
}

start() {
    ebegin "Installing boot-stage dependencies"
    if [ ! -x /usr/local/bin/cpm ]; then
        eerror "FATAL: cpm binary not found at /usr/local/bin/cpm. Image is corrupted."
        return 1
    fi
    /usr/local/bin/cpm install-dependencies --stage boot
    eend $?
}
